# systemd/user/buckle.service
[Unit]
Description=Buckling spring simulator

[Service]
Type=simple
ExecStart=/usr/bin/buckle -f -s 30

[Install]
WantedBy=user-graphical.target

