# systemd/user/urxvtd.service
[Unit]
Description=Urxvt daemon

[Service]
Type=simple
ExecStart=/usr/bin/urxvtd 

[Install]
WantedBy=user-graphical.target

