# systemd/user/xorg.service
[Unit]
Description=Xorg server at display 

[Service]
Type=simple
ExecStart=startx
